# conf/etc/service/nginx/run
#!/bin/sh

mkdir -p /var/log/nginx /var/lib/nginx/logs
exec /usr/sbin/nginx -c /etc/nginx/nginx.conf
